# debian/init4boot-client.postinst
#!/bin/sh

update-rc.d init4boot-client defaults
